# host_ops/qwen36-warning-forwarder.service
[Unit]
# qwen36 §7.5 model_emit_warning forwarder.
#
# Tails `docker logs -f qwen36` and writes structured JSON Lines
# records to /var/log/qwen36/warnings.jsonl whenever the §7.5 patch
# emits ``model_emit_warning kind=tool_call_in_reasoning ...``.
#
# Operational dependencies:
#   * docker.service must be up (the script invokes `docker logs`).
#   * The user the unit runs as must have docker socket access, i.e.
#     be a member of the ``docker`` group, or this unit must run as
#     root. ``User=root`` is the simplest and is what we ship; switch
#     to a docker-group user if your security model requires it.
Description=Qwen3.6 §7.5 model_emit_warning forwarder
Documentation=file:///usr/local/bin/qwen36_warning_forwarder.py
After=docker.service
Requires=docker.service

[Service]
Type=simple
# Run as root for the simplest docker-socket access. To run as a
# non-root user, comment User=root and add the user to the ``docker``
# group, then uncomment User=<your-user>.
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/qwen36_warning_forwarder.py \
    --container qwen36 \
    --output /var/log/qwen36/warnings.jsonl \
    --state  /var/lib/qwen36/forwarder_state.json \
    --no-stdout
Restart=on-failure
RestartSec=5
# 30s grace for SIGTERM cleanup (state flush + output file close).
TimeoutStopSec=30
KillSignal=SIGTERM
# SIGHUP triggers in-place log rotation (close + reopen output file).
# logrotate's ``postrotate`` should send SIGHUP to ``$MAINPID``.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qwen36-warnings

# Light hardening — the script only reads docker logs and writes two
# directories. NoNewPrivileges + a minimal ReadWritePaths allowlist
# keeps blast radius low without breaking docker-socket access.
NoNewPrivileges=true
ReadWritePaths=/var/log/qwen36 /var/lib/qwen36
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
